# deploy/systemd/mcp-heartbeat-collect.service
[Unit]
Description=MCP Heartbeat Collector (writes hb-status.json for all 8 agents)
Documentation=https://github.com/LabDoctorM/mcp-tools
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /root/LabDoctorM/projects/mcp-tools/bin/heartbeat-collect.py
WorkingDirectory=/root/LabDoctorM/projects/mcp-tools
# Read-only except writing each agent's hb-status.json (designated file).
MemoryMax=128M
StandardOutput=journal
StandardError=journal
