# deploy/systemd/mcp-memory.service
[Unit]
Description=MCP Memory Server (labsearch wrapper, read-only)
Documentation=https://github.com/LabDoctorM/mcp-tools
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/LabDoctorM/projects/mcp-tools/bin/memory-server.py
WorkingDirectory=/root/LabDoctorM/projects/mcp-tools
Environment=MCP_TRANSPORT=http
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8087
Restart=always
RestartSec=5
MemoryMax=256M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
