# deploy/systemd/mcp-shell.service
[Unit]
Description=MCP Shell Server (FULL/unrestricted, audited, resource-capped)
After=network.target
Documentation=file:///root/LabDoctorM/projects/mcp-tools/bin/shell_server.py

[Service]
Type=simple
WorkingDirectory=/root/LabDoctorM
ExecStart=/usr/bin/python3 /root/LabDoctorM/projects/mcp-tools/bin/shell_server.py
Environment=MCP_TRANSPORT=http
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8089
Restart=always
RestartSec=5
User=root

# Base hardening only (parity with native exec: network + filesystem writes allowed)
MemoryMax=512M
NoNewPrivileges=true
PrivateTmp=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
