# deploy/systemd/mcp-memory.service
[Unit]
Description=MCP Memory Server (stateful FAISS/ONNX facade)
Documentation=https://github.com/LabDoctorM/mcp-tools
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/LabDoctorM/projects/mcp-tools/bin/memory-server.py
WorkingDirectory=/root/LabDoctorM/projects/mcp-tools
Environment=MCP_TRANSPORT=http
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8087
Environment=STARTUP_BRIEF_PORT=8093
Environment=MCP_BACKEND=hybrid
# Надёжный lifecycle (EBSW 2026-07-11)
Restart=always
RestartSec=5
TimeoutStartSec=60          # FAISS 76MB + meta 32MB грузятся при старте
TimeoutStopSec=15           # graceful shutdown
MemoryMax=2G              # in-memory FAISS + meta + keyword-index + cache
TasksMax=64
# Watchdog НЕ ставим: FastMCP/uvicorn не шлёт sd_notify → systemd убил бы сервер.
# Внутренний _watchdog в коде (60s) перезагружает индекс при падении Штрейкбрехера.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
