# deploy/systemd/mcp-heartbeat-collect.timer
[Unit]
Description=Daily MCP Heartbeat Collector trigger
Documentation=https://github.com/LabDoctorM/mcp-tools
; gatekeeper lease: rk-e64fa1592a77 (daily timer 01:07, as_root bypass — writer under /root)

[Timer]
# Once per day, staggered off midnight (after cron-heavy infra settles).
OnCalendar=*-*-* 01:07:00
Persistent=true
RandomizedDelaySec=300
Unit=mcp-heartbeat-collect.service

[Install]
WantedBy=timers.target
