# mcp-gatekeeper/systemd/mcp-gatekeeper.service
[Unit]
Description=MCP Gatekeeper (привратник портов/таймеров агентов)
Documentation=file:///root/LabDoctorM/projects/mcp-tools/mcp-gatekeeper/docs/CONTRACT.md
After=network.target
# StartLimit* ограничивают crash-loop (это [Unit]-секция, не [Service])
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/LabDoctorM/projects/mcp-tools/mcp-gatekeeper/bin/mcp-gatekeeper-server.py
WorkingDirectory=/root/LabDoctorM/projects/mcp-tools/mcp-gatekeeper
Environment=MCP_TRANSPORT=http
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8888
Environment=GATEKEEPER_POLICY=/root/LabDoctorM/projects/mcp-tools/mcp-gatekeeper/policies/policy_v1.yaml
Environment=GATEKEEPER_DATA=/root/LabDoctorM/projects/mcp-tools/mcp-gatekeeper/data

# Надёжный lifecycle (EBSW-совместимо). Зависание/смерть процесса
# ловится событийно через Restart (без polling watchdog): агент ставит
# таймаут 5с на вызов register_* и при таймауте пишет GATEKEEPER_TIMEOUT.
Restart=on-failure
RestartSec=5
# Периодический опрос (watchdog) удалён: polling-модель избыточна (см. docs/CONTRACT.md)
# Config validation fail-fast: если политика невалидна — быстрый выход,
# systemd сделает on-failure рестарт, упрётся в StartLimitBurst и остановится.
TimeoutStartSec=30
TimeoutStopSec=15
MemoryMax=128M

# --- Hardening (Фаза 5 ADR-0055): RBAC/защита конфигурации (hole 2/8) ---
# Файлы policy/leases — root:gatekeeper 640 (см. set-rbac.sh); процесс не может
# их править/удалить. ProtectSystem=strict делает ФС read-only для процесса,
# кроме явно разрешённых путей записи (data/ для leases.json и журнала).
ProtectSystem=strict
# data/ должен быть доступен на запись (lease-состояние + JSONL-журнал):
ReadWritePaths=/root/LabDoctorM/projects/mcp-tools/mcp-gatekeeper/data
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
# Лишний сетевой доступ не нужен, но слушать 127.0.0.1:8888 надо — это разрешено.
# (RestrictAddressFamilies не задаём, чтобы не сломать bind на localhost.)

# Least-privilege (опционально): ресурсы выдаются «от ограниченного юзера»,
# не от root. Раскомментируйте и создайте юзера:
#   useradd --system --no-create-home --shell /usr/sbin/nologin mcp-gatekeeper
#   chown -R mcp-gatekeeper:gatekeeper data
# User=mcp-gatekeeper
# Тогда ReadWritePaths выше должен указывать на data/, доступный этому юзеру.

# Log rate-limit реализован в коде (защита диска).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
