# memory-gateway/systemd/alm-sync-rebuild.service
[Unit]
Description=AnythingLLM full reindex (rebuild) + lexical convergence
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/LabDoctorM/.ops/shared/anythingllm-sync/run-sync.sh rebuild
User=root
Nice=10
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
