# memory-gateway/systemd/alm-sync-incremental.timer
[Unit]
Description=AnythingLLM incremental sync every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=false
Unit=alm-sync-incremental.service

[Install]
WantedBy=timers.target
